# crons/submgr.crontab
MAILTO=boropeza@peak6.com,jflude@peak6.com,jrobertson@peak6.com,mvanderpluym@peak6.com

05 01 * * mon-fri sleep `shuf -i 60-180 -n 1` ; exec submgr/latest/scripts/cron.submgr.start

45 23 * * mon-fri submgr/latest/scripts/cron.submgr.stop ; exec submgr/latest/scripts/cron.submgr.roll
